# pi-flirc-roon-volume/soundstream-flirc-roon-volume.service
[Unit]
Description=SoundStream FLIRC -> Roon volume listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/soundstream-flirc-roon-volume.env
ExecStart=/opt/soundstream-flirc-roon-volume/venv/bin/python /opt/soundstream-flirc-roon-volume/flirc_roon_volume.py
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target


